# metamix-supervisor.conf
[program:metamix]
command = gunicorn run_app:app -b localhost:8000
directory = /app/metamix/
user = metamix-user